# listener/upstart-pull.conf
description "kdd2be-pull"

start on (filesystem)
stop on runlevel [016]

respawn
console log
setuid kdd2
setgid kdd2
chdir /home/kdd2/tester

exec ./pull.py
